# frontend/app/config/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:streamlit]
command=python APP_wrap.py